# 10-block.rules
SUBSYSTEM=="block", RUN+="/sbin/rebuildfstab"
